# cron.hourly
#!/bin/sh

/opt/fffd-dns/dns-nodenames.py > /opt/fffd-dns/zones/db.nodes.fffd.eu
rndc reload > /dev/null

exit 0
